# conf/etc/services.d/01-omnidb/run
#!/bin/ash
NAME="OMNIDB" 
COMMAND="/usr/local/bin/omnidb-launcher" 
WAIT="/tmp/omnidb"
exec service-handler "$NAME" "$COMMAND" "$WAIT"